# build/package/rpm/alertmanager_matrix.service.in
[Unit]
Description=Alerts from Alertmanager to Matrix
After=network-online.target

[Service]
Restart=always
DynamicUser=yes
EnvironmentFile=@DEFAULTDIR@/alertmanager_matrix
ExecStart=@BINDIR@/alertmanager_matrix

[Install]
WantedBy=multi-user.target
